# etc/systemd/system/slubstart.service
[Unit]
Description=SLUB Startup

[Service]
Type=simple
ExecStart=/opt/slub/startup.sh

[Install]
WantedBy=multi-user.target display-manager.service
